# cinder/files/cinder-volume-usage-audit.cron
#!/bin/bash

export LANG=en_US.UTF-8
export PYTHONIOENCODING=utf-8

/usr/bin/python -u /usr/bin/cinder-volume-usage-audit 1>> /var/log/cinder/cinder_audit_hourly__`date +\%Y\%m\%d`.log 2>&1

EXITVALUE=$?
if [ $EXITVALUE != 0 ]; then
    /usr/bin/logger -t cinder-volume-usage-audit "ALERT exited abnormally with [$EXITVALUE]"
fi
exit 0
